# cookbooks/fefiles/files/default/etc_cron_bin/cron.weekly/anteater-fe-weekly
#!/bin/bash

# Placed by chef

/usr/bin/find /home/ant/.mythtv/remotecache -type f -atime +90 -print0 | xargs -0 /bin/rm -f
