# cookbooks/befiles/files/default/etc_cron_bin/cron.monthly/mythtv-anteater-monthly
#!/bin/bash

# Placed by chef

/usr/local/bin/clonetxt | /usr/bin/mailx -s "DVR Clone-Backup reminder" linuxtrader@gmail.com

